# deploy/whatsapp-mcp.service
[Unit]
Description=WhatsApp MCP persistent daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/Code/whatsapp-mcp/deploy/whatsapp-mcp-serve
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=PATH=%h/.local/share/mise/shims:%h/.mise/shims:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
