# deploy/whatsapp-mcp.service
[Unit]
Description=WhatsApp MCP persistent daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/Code/whatsapp-mcp/deploy/whatsapp-mcp-serve
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
