# deploy/whatsapp-mcp.service
[Unit]
Description=WhatsApp MCP persistent daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/Code/whatsapp-mcp/deploy/whatsapp-mcp-serve
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
